# environments/debian.postinst
#!/bin/sh
gpasswd -a testbench docker
systemctl enable NetworkManager
